# systemd/minecraft-backup.service
[Unit]
Description=Restic backup of the Minecraft server
Documentation=https://restic.readthedocs.io/en/stable/
Requires=minecraft.service
After=minecraft.service network-online.target
Wants=network-online.target
ConditionPathIsDirectory=/var/lib/minecraft/data
ConditionPathIsDirectory=/var/lib/minecraft/backup
ConditionPathExists=/var/lib/minecraft/data/rcon.yaml

[Service]
Type=oneshot
TimeoutStartSec=infinity
ExecSearchPath=/home/linuxbrew/.linuxbrew/bin:/home/linuxbrew/.linuxbrew/sbin
ExecSearchPath=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin

# Commands run serially. A failure skips the remaining commands and runs ExecStopPost.
ExecStart=rcon-cli \
  --config=/var/lib/minecraft/data/rcon.yaml \
  save-off

ExecStart=rcon-cli \
  --config=/var/lib/minecraft/data/rcon.yaml \
  save-all flush

ExecStart=restic \
  --repo=/var/lib/minecraft/backup \
  --insecure-no-password \
  backup \
  /var/lib/minecraft/data

ExecStart=restic \
  --repo=/var/lib/minecraft/backup \
  --insecure-no-password \
  forget \
  --keep-last=20 \
  --prune

# ExecStopPost also runs when backup or pruning fails.
ExecStopPost=-rcon-cli \
  --config=/var/lib/minecraft/data/rcon.yaml \
  save-on
